# scripts/cron/traders_family_etl.cron
# Run the Traders Family ETL every day at 08:00 server local time.
# Optional: prepend `CRON_TZ=Asia/Jakarta` above this line if the server timezone is not WIB.
CRON_TZ=Asia/Jakarta
00 08 * * * /bin/bash {{PROJECT_ROOT}}/scripts/run_scheduled_etl.sh
